# systemd/edudip-import.service
[Unit]
Description=Import edudip webinar recordings

[Service]
Type=simple
ExecStart=/opt/edudip/import
User=edudip
